# deploy/supervisor.conf
; Supervisor configuration for Blog FastAPI backend
; Place this file in /etc/supervisor/conf.d/blog.conf

[program:blog-api]
command=/var/www/blog/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
directory=/var/www/blog/backend
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/blog-api.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=
    PATH="/var/www/blog/backend/.venv/bin:%(ENV_PATH)s",
    PYTHONUNBUFFERED="1"

; Optional: Use systemd instead of supervisor
; Create /etc/systemd/system/blog-api.service with:
;
; [Unit]
; Description=Blog FastAPI Backend
; After=network.target
;
; [Service]
; User=www-data
; Group=www-data
; WorkingDirectory=/var/www/blog/backend
; Environment="PATH=/var/www/blog/backend/.venv/bin"
; ExecStart=/var/www/blog/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
; Restart=always
;
; [Install]
; WantedBy=multi-user.target
